# bin/e-Paper_clean.service
[Unit]
Description=Run e-Paper clean on shutdown and reboot
DefaultDependencies=no
Before=shutdown.target
Before=reboot.target

[Service]
Type=oneshot
ExecStart=/bin/bash /root/2.13-Ink-screen-clock/bin/clean.sh
RemainAfterExit=yes
TimeoutSec=5

[Install]
WantedBy=shutdown.target
WantedBy=reboot.target
