# bin/raspi_e-Paper.service
[Unit]
Description=Raspi e-Paper Service
After=network.target

[Service]
Type=forking
ExecStart=/root/2.13-Ink-screen-clock/bin/start.sh
ExecStop=/bin/bash -c "pkill -f /root/2.13-Ink-screen-clock/bin/main.py"
WorkingDirectory=/root/2.13-Ink-screen-clock/bin
Restart=always
User=root

[Install]
WantedBy=multi-user.target