# bin/raspi_e-Paper.service
[Unit]
Description=Raspi e-Paper Service (C Version)
After=network.target

[Service]
Type=forking
ExecStart=/root/2.13-Ink-screen-clock/bin/start.sh
ExecStop=/bin/bash -c "pkill -f epd_clock; pkill -f weather.py"
WorkingDirectory=/root/2.13-Ink-screen-clock/bin
Restart=always
User=root

[Install]
WantedBy=multi-user.target
