# bin/raspi_e-Paper.service
[Unit]
Description=Raspi e-Paper Service
After=network.target

[Service]
Type=forking
# 等待网络完全启动
ExecStartPre=/bin/bash -c 'until ping -c1 baidu.com &>/dev/null; do sleep 1; done'
ExecStart=/root/2.13-Ink-screen-clock/bin/start.sh
ExecStop=/usr/bin/screen -S main_screen -X quit
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
